# jokes-docker/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:jokes]
priority=1
command=java -jar /jokes.jar -Dspring.profiles.active=test
stderr_logfile = /tmp/jokes-stderr.log
stdout_logfile = /tmp/jokes-stdout.log
